# supervisor/motion-detector.conf
[program:motion-detector]
command=/usr/sbin/motion-detection.py
autostart=true
autorestart=true
stderr_logfile=/var/log/network-monitor.log
