# supervisor/network-monitor.conf
[program:network-monitor]
command=/usr/sbin/network-monitor.py
autostart=true
autorestart=true
stderr_logfile=/var/log/network-monitor.log
